# config/supervisord.ini
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
user=root
numprocs=1
autostart=true
autorestart=false
startsecs=0
loglevel=trace
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm]
command=php-fpm -R -F -y /usr/local/etc/php/php-fpm.conf
user=root
numprocs=1
autostart=true
autorestart=false
startsecs=0
loglevel=trace
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
